# packages/macos-virt/macos-virt-0.2.0.tar.gz/macos-virt-0.2.0/macos_virt/service/macos-virt-service.service
[Unit]

Description=Macos-Virt Agent.

[Service]

ExecStart=/bin/python3 -u /usr/sbin/macos-virt-service.py

[Install]

WantedBy=multi-user.target